# ptt-script/push-to-talk.service
[Unit]
Description=Push-to-Talk Service
After=pipewire.service pulseaudio.service

[Service]
Type=simple
ExecStart=%h/.local/bin/push-to-talk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
